# mac/pkg/scripts/postinstall
#!/bin/bash
set -euo pipefail

# --- Resolve the active console user (for Desktop alias + starting tray) ---
CURRENT_USER="$(/usr/sbin/scutil <<< 'show State:/Users/ConsoleUser' | awk '/Name :/ && $3 != "loginwindow" {print $3}' || true)"
if [[ -z "${CURRENT_USER:-}" || "${CURRENT_USER}" == "root" ]]; then
  CURRENT_USER="$(stat -f%Su /dev/console || true)"
fi

APP="/Applications/NOVA.app"
TRAY_APP="/Applications/Nova Tray.app/Contents/MacOS/NovaTray"

# --- Install LaunchAgent to auto-start Nova Tray at login (ALWAYS do this, even on CI) ---
SCRIPT_DIR="$(cd "$(dirname "$0")" && pwd)"
PLIST_SRC="${SCRIPT_DIR}/com.novaai.tray.plist"
PLIST_DST="/Library/LaunchAgents/com.novaai.tray.plist"
/bin/mkdir -p "/Library/LaunchAgents"
/bin/cp -f "$PLIST_SRC" "$PLIST_DST"
/usr/sbin/chown root:wheel "$PLIST_DST"
/bin/chmod 644 "$PLIST_DST"

# --- If running on CI, stop here: no GUI/session actions on headless runners ---
if [[ "${GITHUB_ACTIONS:-}" == "true" || "${CI:-}" == "true" ]]; then
  exit 0
fi

# --- Create a Desktop alias to NOVA.app for the console user (best-effort) ---
if [[ -n "${CURRENT_USER:-}" && "${CURRENT_USER}" != "root" ]]; then
  DESK="/Users/${CURRENT_USER}/Desktop"
  /bin/mkdir -p "$DESK"

  # Try a Finder alias first (requires GUI); fall back to a symlink.
  if pgrep -x Finder >/dev/null 2>&1; then
    /usr/bin/osascript <<'OSA' || true
tell application "Finder"
  try
    set d to (path to desktop folder) as alias
    set targetApp to POSIX file "/Applications/NOVA.app" as alias
    if not (exists file "NOVA.app" of d) then
      make new alias file at d to targetApp with properties {name:"NOVA.app"}
    end if
  end try
end tell
OSA
  fi

  [[ -e "$DESK/NOVA.app" ]] || /bin/ln -sfn "$APP" "$DESK/NOVA.app" || true
  /usr/sbin/chown -h "$CURRENT_USER":staff "$DESK/NOVA.app" 2>/dev/null || true
fi

# --- Start (or restart) the tray in the current GUI session (best-effort) ---
if [[ -n "${CURRENT_USER:-}" && -x "$TRAY_APP" ]]; then
  UID="$(/usr/bin/id -u "$CURRENT_USER" 2>/dev/null || true)"
  if [[ -n "${UID:-}" ]]; then
    /bin/launchctl bootout   "gui/${UID}/com.novaai.tray" >/dev/null 2>&1 || true
    /bin/launchctl bootstrap "gui/${UID}" "$PLIST_DST"     >/dev/null 2>&1 || true
    /bin/launchctl enable    "gui/${UID}/com.novaai.tray"  >/dev/null 2>&1 || true
    /bin/launchctl kickstart -k "gui/${UID}/com.novaai.tray" >/dev/null 2>&1 || true
  fi
fi

exit 0
